# quadlets/homeserver-mosquitto.container
[Unit]
Description=Mosquitto MQTT Broker for Home Server
After=homeserver-network.service
Requires=homeserver-network.service

[Container]
Image=docker.io/eclipse-mosquitto:latest
ContainerName=homeserver-mosquitto
AutoUpdate=registry

# Network - CRITICAL: Static IP!
Network=homeserver.network
IP=172.20.0.10
PublishPort=1883:1883
PublishPort=9001:9001

# Volumes - UPDATE paths for your system
Volume=%h/homeserver/configs/mosquitto:/mosquitto/config:Z
Volume=%h/homeserver/data/mosquitto/data:/mosquitto/data:Z
Volume=%h/homeserver/data/mosquitto/log:/mosquitto/log:Z

# Environment
Environment=TZ=America/New_York

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
