# quadlets/homeserver-zigbee2mqtt.container
[Unit]
Description=Zigbee2MQTT Bridge for Home Server
After=homeserver-mosquitto.service
Requires=homeserver-mosquitto.service

[Container]
Image=docker.io/koenkk/zigbee2mqtt:latest
ContainerName=homeserver-zigbee2mqtt
AutoUpdate=registry

# Network - CRITICAL: Static IP!
Network=homeserver.network
IP=172.20.0.40
PublishPort=8080:8080

# Volumes - UPDATE paths for your system
Volume=%h/homeserver/configs/zigbee2mqtt:/app/data:Z
Volume=/run/udev:/run/udev:ro

# Device passthrough for Zigbee coordinator - UPDATE for your device!
# Find with: ls -la /dev/serial/by-id/
AddDevice=/dev/serial/by-id/usb-Silicon_Labs_slae.sh_cc2652rb_stick_-_slaesh_s_iot_stuff_00_12_4B_00_25_9B_BF_73-if00-port0:rwm

# Environment
Environment=TZ=America/New_York

# Run with host user groups (for dialout access)
GroupAdd=keep-groups

# Security - disable SELinux for device access
SecurityLabelDisable=true
PodmanArgs=--privileged

# Health Check
HealthCmd=wget --no-verbose --tries=1 --spider http://localhost:8080 || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
