# quadlets/homeserver-homeassistant.container
[Unit]
Description=Home Assistant for Home Server
After=homeserver-mosquitto.service
Wants=homeserver-mosquitto.service homeserver-zigbee2mqtt.service

[Container]
Image=ghcr.io/home-assistant/home-assistant:stable
ContainerName=homeserver-homeassistant
AutoUpdate=registry

# Network - CRITICAL: Static IP!
Network=homeserver.network
IP=172.20.0.30
PublishPort=8123:8123

# Volumes - UPDATE paths for your system
Volume=%h/homeserver/configs/homeassistant:/config:Z
Volume=/etc/localtime:/etc/localtime:ro
Volume=/run/dbus:/run/dbus:ro

# Environment
Environment=TZ=America/New_York

# Install Frigate dependency on startup (CRITICAL for Frigate integration!)
Exec=/bin/sh -c "echo nameserver 8.8.8.8 > /etc/resolv.conf && pip3 install -q hass-web-proxy-lib==0.0.7 2>/dev/null; exec /init"

# Health Check
HealthCmd=curl --fail http://localhost:8123 || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=5
HealthStartPeriod=120s

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
