# quadlets/homeserver-frigate.container
[Unit]
Description=Frigate NVR for Home Server
After=homeserver-mosquitto.service
Wants=homeserver-mosquitto.service

[Container]
Image=ghcr.io/blakeblackshear/frigate:stable
ContainerName=homeserver-frigate
AutoUpdate=registry

# Network - CRITICAL: Static IP!
Network=homeserver.network
IP=172.20.0.20
PublishPort=5000:5000
PublishPort=8554:8554
PublishPort=8555:8555/tcp
PublishPort=8555:8555/udp

# Volumes - UPDATE paths for your system
Volume=%h/homeserver/configs/frigate:/config:Z
Volume=%h/homeserver/data/media/frigate:/media/frigate:Z
Volume=/etc/localtime:/etc/localtime:ro

# Shared memory for faster processing
ShmSize=512mb

# Environment
Environment=TZ=America/New_York

# Capabilities
AddCapability=CAP_PERFMON

# Health Check
HealthCmd=curl --fail http://localhost:5000/api/version || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
