# quadlets/homeserver-mosquitto.container
[Unit]
Description=Mosquitto MQTT Broker for Home Server
After=homeserver-network.service
Requires=homeserver-network.service

[Container]
Image=docker.io/eclipse-mosquitto:latest
ContainerName=homeserver-mosquitto
AutoUpdate=registry

# Network
Network=homeserver.network
PublishPort=1883:1883
PublishPort=9001:9001

# Volumes
Volume=%h/Documents/GitHub/goofaround/configs/mosquitto:/mosquitto/config:Z
Volume=%h/Documents/GitHub/goofaround/data/mosquitto/data:/mosquitto/data:Z
Volume=%h/Documents/GitHub/goofaround/data/mosquitto/log:/mosquitto/log:Z

# Environment
Environment=TZ=America/New_York

# Security
User=1883
SecurityLabelDisable=false

# Health Check
HealthCmd=mosquitto_sub -t '$$SYS/#' -C 1 || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target


